# root/etc/services.d/guacd/run
#!/usr/bin/with-contenv sh

echo "Starting guacamole guacd..."
s6-setuidgid root guacd -f
